# 14.x/x86_64/tcz/oracle-oci8/src/etc/profile.d/oracle-oci8.sh
export ORACLE_CLIENT_HOME="/usr/local/oracle-oci8"
export LD_LIBRARY_PATH=${ORACLE_CLIENT_HOME}:$LD_LIBRARY_PATH
export NLS_LANG="AMERICAN_AMERICA.AL32UTF8"
export PATH=${ORACLE_CLIENT_HOME}:$PATH
